# Definitions/main.proxy.caching/docker/squid/1.0.0/proto.squid.service
[Unit]
Description=Squid running on Docker
After=docker.service
Requires=docker.service

[Service]
Restart=always
ExecStart={{DOCKER.DOCKER_HOME}}/Definitions/Squid/start.sh

[Install]
WantedBy=multi-user.target